# ops/web/files/crontab
*/1 * * * * sh /project/ops/web/sync-git.sh
